# deploy/mlapi.service
[Unit]
Description=OKR Project ML API Server
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/okr-project
Environment=PATH=/home/ubuntu/okr-project/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=PYTHONPATH=/home/ubuntu/okr-project
Environment=FLASK_ENV=production
Environment=FLASK_APP=api/app.py
Environment=PORT=5001
Environment=HOST=0.0.0.0

# Start the API server using the virtual environment
ExecStart=/home/ubuntu/okr-project/venv/bin/python api/app.py

# Restart configuration
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/okr-project

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=okr-api

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
